# scripts/whenimfree.service
[Unit]
Description=WhenImFree
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/whenimfree
Environment=PYTHONPATH=/home/ubuntu/whenimfree/src
ExecStart=/home/ubuntu/.local/bin/uv run uvicorn whenimfree.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
